# hemanth-bot.service
[Unit]
Description=вҡЎрқ—ӣрқ—ҳрқ— рқ—”рқ—Ўрқ—§рқ—ӣвҡЎ Mirror-Leech Bot Service
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/mirrorbot
Environment="PATH=/root/mirrorbot/venv/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
ExecStart=/root/mirrorbot/venv/bin/python3 -m bot
Restart=always
RestartSec=10
KillSignal=SIGINT

[Install]
WantedBy=multi-user.target
